# src/init_script
#!/bin/sh /etc/rc.common

START=1000
STOP=1015

restart() {
  echo restart
  opkg update
  opkg install python -d ram
  opkg install vim -d ram
  opkg install tmux -d ram
}
